# src/swift-cfg/res/swift-cfg.service
[Unit]
Description=Swift Desktop Configuration Service
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/swift-cfg
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target